# api/nockchain-jammer-api.service
[Unit]
Description=Nockchain Jammer API
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nockchain-jammer-api
Environment=API_KEY=CHANGE_ME
Environment=SCRIPT_PATH=/usr/local/bin/make-jam.sh
Environment=HTML_ROOT=/usr/share/nginx/html
Environment=JAMS_DIR=/usr/share/nginx/html/jams
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
